# docs/ops/services/launchd_p99_ops_loop_guarded_v1.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.peaktrade.p99-ops-loop-guarded</string>

  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>300</integer>

  <key>WorkingDirectory</key>
  <string>/path/to/Peak_Trade</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>cd "$HOME/Peak_Trade" &amp;&amp; ./scripts/ops/p99_ops_loop_launchd_guard_v1.sh</string>
  </array>

  <key>StandardOutPath</key>
  <string>out/ops/online_readiness_supervisor/LAUNCHD_P99_OPS_LOOP_GUARDED.out.log</string>
  <key>StandardErrorPath</key>
  <string>out/ops/online_readiness_supervisor/LAUNCHD_P99_OPS_LOOP_GUARDED.err.log</string>
</dict>
</plist>
